# misc/packages/compton/compton.service
[Unit]
Description=Compositing manager
After=graphical-session.target

[Service]
ExecStart=/usr/bin/compton --backend glx
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
