# misc/packages/compton/compton.service
[Unit]
Description=Compositing manager
After=xorg.target

[Service]
ExecStart=/usr/bin/compton
Environment=DISPLAY=:0

[Install]
WantedBy=wm.target
