# misc/packages/compton/compton.service
[Unit]
Description=Compositing manager
After=xorg.target dbus.service
Requires=dbus.service

[Service]
ExecStart=/usr/bin/compton
Environment=DISPLAY=:0
Restart=always

[Install]
WantedBy=wm.target
