# misc/packages/xmonad/plasma-xmonad.service
[Unit]
Description=Plasma with XMonad window manager
Before=plasma-workspace.target

[Service]
ExecStart=%h/.local/bin/xmonad
Slice=session.slice
Restart=on-failure

[Install]
WantedBy=plasma-workspace.target
