# misc/packages/homelink/homelink-monitor.service
[Unit]
Description=HomeLink Connectivity Monitor
After=homelink.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelink-monitor.sh
