# misc/packages/homelink/homelink-monitor.service
[Unit]
Description=HomeLink Connectivity Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/homelink-monitor.sh

[Install]
WantedBy=multi-user.target
