# misc/packages/homelink/install
#!/usr/bin/env bash

script_dir="$(cd "$(dirname "${BASH_SOURCE[0]}")" && pwd)"

echo "Copy homelink.service -> /etc/systemd/system/homelink.service"
sudo cp "$script_dir/homelink.service" /etc/systemd/system/homelink.service
sudo chmod 644 /etc/systemd/system/homelink.service

echo "Copy homelink-monitor.service -> /etc/systemd/system/homelink-monitor.service"
sudo cp "$script_dir/homelink-monitor.service" /etc/systemd/system/homelink-monitor.service
sudo chmod 644 /etc/systemd/system/homelink-monitor.service

echo "Copy homelink script -> /usr/local/bin/homelink"
sudo cp "$script_dir/homelink" /usr/local/bin/homelink
sudo chmod 755 /usr/local/bin/homelink

echo "Reload systemd and enable monitor service..."
sudo systemctl daemon-reload
sudo systemctl enable homelink-monitor.service

echo "Start monitor service..."
sudo systemctl start homelink-monitor.service
