# misc/packages/nb/nb-sync.timer
[Unit]
Description=Run nb notebook sync every 10 minutes
After=graphical-session.target

[Timer]
Unit=nb-sync.service
OnStartupSec=5min
OnUnitInactiveSec=10m
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
