# misc/packages/nb/nb-sync.service
[Unit]
Description=Run nb notebook sync
Wants=nb-sync.timer

[Service]
ExecStart=%h/.local/bin/nb sync --all --no-color
Environment=NB_DIR=%h/.mcf/misc/packages/nb/notebooks TERM=xterm
