# misc/packages/nb/nb-sync.timer
[Unit]
Description=Run nb notebook sync every 10 minutes
Requires=nb-sync.service

[Timer]
Unit=nb-sync.service
OnBootSec=5min
OnUnitInactiveSec=10m
RandomizedDelaySec=2m
AccuracySec=1s

[Install]
WantedBy=timers.target
